# bin/space_machines.service
[Unit]
Description=Sector67 Access Control
After=multi-user.target

[Service]
Type=forking
PIDfile=#PIDPATH#
ExecStart=#STARTPATH#

[Install]
WantedBy=multi-user.target